# com.local-ai-agents.telegram-bot.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.local-ai-agents.telegram-bot</string>
    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/uv</string>
        <string>run</string>
        <string>python</string>
        <string>-u</string>
        <string>scripts/telegram_bot.py</string>
    </array>
    <key>WorkingDirectory</key>
    <!-- Update this path to match your local checkout -->
    <string>/Users/YOURUSERNAME/GitHub/local-ai-agents</string>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/YOURUSERNAME/GitHub/local-ai-agents/data/telegram-bot.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/YOURUSERNAME/GitHub/local-ai-agents/data/telegram-bot.err</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
